# comando21-script6.sh
#! /bin/sh
### BEGIN INIT INFO
# IFTO - GREDES
# Provides:          CCM
# Documentation:     Lucas Coutinho de Almeida
# Short-Description: 
### END INIT INFO

#This script, executes two commands, comando1 and restart winbind
/bin/comando1
/etc/init.d/winbind restart