# comando4-script19.sh
#! /bin/sh
### BEGIN INIT INFO
# IFTO - GREDES
# Provides:          CCM
# Documentation:     Lucas Coutinho de Almeida
# Short-Description: 
### END INIT INFO

#Remove all the folders and files of read from /home/aluno/* 
rm -Rf /home/aluno/*