# limpaaluno.sh
#! /bin/sh
### BEGIN INIT INFO
# IFTO - GREDES
# Provides:          CCM
# Documentation:     Paulo Vinicius Rodrigues Borges
# Short-Description: remove directory
### END INIT INFO

#!/bin/bash
#remove directory dr froma recursimanevente and without confirmation
rm -rf /home/aluno/*



